# mt7988a/tcp_client.init
#!/bin/sh /etc/rc.common

START=99
STOP=10

USE_PROCD=1

start_service() {
    echo "Starting MT7988A TCP Client..."

    procd_open_instance
    procd_set_param command /root/run_client.sh
    procd_set_param respawn
    procd_close_instance
}

stop_service() {
    echo "Stopping MT7988A TCP Client..."
}
